# crontab
@reboot ~/bin/start
@reboot ~/bin/start-monitor
